# configuration/role-app-go/files/default/goapp.conf
# goapp service script
start on startup
exec /opt/go/bin/goapp
